# doc/systemd/svxlink.timer
[Timer]
OnBootSec=30

[Install]
WantedBy=graphical.target
